# deploy/registry.service
[Unit]
Description=Agent Auth Registry
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent-auth
Group=agent-auth
WorkingDirectory=/var/lib/agent-auth
EnvironmentFile=/etc/agent-auth/registry.env
ExecStart=/opt/agent-auth/venv/bin/agent-auth-registry
Restart=on-failure
RestartSec=3
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
LockPersonality=true
ReadWritePaths=/var/lib/agent-auth

[Install]
WantedBy=multi-user.target
